# server/elder-berry.service
[Unit]
Description=Elder-Berry Saleria Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=lera
Group=psacln
WorkingDirectory=/opt/Elder-Berry
Environment="PATH=/opt/Elder-Berry/.venv/bin:/usr/local/bin:/usr/bin"
Environment="ELDER_BERRY_HOME=/opt/Elder-Berry"
ExecStart=/opt/Elder-Berry/.venv/bin/python scripts/start_saleria.py --mode matrix
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sicherheit
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
